# bin/musl-x64/duckdns.service
[Unit]
Description=Duck DNS Client

[Service]
Type=notify
ExecStart=/usr/bin/ddnsservice
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target